# docker/supervisord.conf
[supervisord]
logfile=/home/printer/klipper_logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10 ; (num of main logfile rotation backups;default 10)
loglevel=info ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true ; (start in foreground if true;default false)
minfds=1024 ; (min. avail startup file descriptors;default 1024)
minprocs=200 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[program:simulavr]
command=bash -c 'sleep 10 && /home/printer/klipper/scripts/avrsim.py -m atmega644 /home/printer/simulavr.elf'
environment=PYTHONPATH=/home/printer/simulavr/build/pysimulavr/
user=printer
process_name=simulavr
directory=/home/printer
autostart=true
autorestart=true
redirect_stderr=true

[program:klipper]
command=/home/printer/klippy-env/bin/python /home/printer/klipper/klippy/klippy.py /home/printer/klipper_config/printer.cfg -l /home/printer/klipper_logs/klippy.log -a /tmp/klippy_uds
user=printer
process_name=klipper
directory=/home/printer
autostart=true
autorestart=true
redirect_stderr=true

[program:moonraker]
command=/home/printer/moonraker-env/bin/python /home/printer/moonraker/moonraker/moonraker.py -l /home/printer/klipper_logs/moonraker.log -c /home/printer/klipper_config/moonraker.conf
user=printer
process_name=moonraker
directory=/home/printer
autostart=true
autorestart=true
redirect_stderr=true